# deploy/systemd/otterlink.target
[Unit]
Description=OtterLink stack
Wants=otterlink-core.service otterlink-gateway.service
After=otterlink-core.service otterlink-gateway.service
