# deploy/systemd/remoteagent.target
[Unit]
Description=remoteagent stack
Wants=remoteagent-core.service remoteagent-gateway.service
After=remoteagent-core.service remoteagent-gateway.service
