# deploy/launchd/com.otterlink.core.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.otterlink.core</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>exec "__ROOT_DIR__/scripts/launchd-core.sh"</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>OTTERLINK_ENV_FILE</key>
      <string>__ENV_FILE__</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__ROOT_DIR__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__ROOT_DIR__/.run/core.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__ROOT_DIR__/.run/core.launchd.log</string>
  </dict>
</plist>
